# pi-assistant.service
[Unit]
Description=Pi Assistant
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=pi-assistant
Restart=on-failure
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
